# supervisor/ghost.conf

[program:ghost]
command = node /ghost/index.js
directory = /ghost
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/ghost.log
stderr_logfile = /var/log/supervisor/ghost_err.log
environment = NODE_ENV="production"
